# user/openlumi-19.07/files/files/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

echo 30000 >  /proc/sys/vm/min_free_kbytes
echo 2 >/proc/sys/kernel/printk
echo -e "\033[9;0]" > /dev/tty0
echo 0 > /sys/bus/platform/drivers/imx_busfreq/soc\:busfreq/enable

if [ -f /etc/asound.state ]; then
  /usr/sbin/alsactl -f /etc/asound.state restore
fi

for i in 1 0; do
  for c in red green blue; do
    echo $i > /sys/class/leds/$c/brightness
  done
done

exit 0
