# user/openlumi-21.02.0-rc/files/target/linux/imx6/base-files/etc/uci-defaults/50_disable-removing-wdev
. /lib/functions.sh

sed -i 's/iw dev "$wdev" del/# iw dev "$wdev" del/' /lib/netifd/wireless/mac80211.sh

exit 0
