# user/openlumi-19.07.7/files/files/etc/uci-defaults/70_disable_wan_firewall
uci -q batch <<-EOF >/dev/null
	set firewall.@zone[0].input='ACCEPT'
	set firewall.@zone[1].input='ACCEPT'
	commit network
EOF

exit 0
